# files/uci-defaults/01-simet1-cleanup
#!/bin/sh

rm -f /etc/config/simetbox_hash_configured /etc/config/hash 2>/dev/null
exit 0
