# files/uci-defaults/91-offline
#!/bin/sh
uci set system.status=status
uci set system.status.online=0
exit 1